# build/cdcm.service
[Unit]
Description=Credentialed Data Collection Module
Requires=network.target
After=network.target

[Service]
Environment=LD_LIBRARY_PATH=${LD_LIBRARY_PATH}:/usr/share/cdcm/lib
ExecStart= /usr/bin/cdcm_supervisor
ExecStartPre = -/usr/bin/rm /dev/shm/sessions_lock
Restart=always

[Install]
WantedBy=multi-user.target
